# data/etc/rc.d/mdnshubd.OpenBSD.in
#!/bin/sh

daemon="@CMAKE_INSTALL_PREFIX@/sbin/mdnshubd"
daemon_user="_mdnshubd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	install -d -o ${daemon_user} /var/run/mdnshubd
}

rc_cmd $1

